# bty-media/rootfs/server/etc/systemd/system/systemd-networkd-wait-online.service.d/any.conf
# bty server appliance override: pass ``--any`` so wait-online is
# satisfied as soon as ANY managed interface is routable, instead of
# the default "wait for all of them". A PXE server with a dead
# segment NIC (no cable, no peer, etc.) shouldn't block boot for the
# full 120s default timeout - bty-web doesn't need that NIC online,
# only the mgmt one.

[Service]
ExecStart=
ExecStart=/lib/systemd/systemd-networkd-wait-online --any --timeout=30
