# bty-media/live-build/config/includes.chroot/etc/systemd/system/bty-tui-on-tty1.service
[Unit]
Description=bty-tui interactive flash on tty1
Documentation=https://github.com/safl/bty
# Fire whenever ``bty.mode=interactive`` is on the kernel cmdline.
# Two paths set it:
#   - PXE-tui: server's ``ipxe_tui.j2`` (``boot_policy=tui``)
#   - USB-local: the usb-iso live-build bake sets it directly via
#     ``--bootappend-live`` (M19 phase 2).
# ``ipxe_flash.j2`` (``boot_policy=flash``) does NOT set it, so the
# auto-flash boots fall through to ``bty-flash-on-boot.service``.
ConditionKernelCommandLine=bty.mode=interactive
# Take over tty1: the live env's auto-login agetty is the default
# session there; we want to BE the session, not run inside it.
# ``Before=getty@tty1.service`` (combined with ``Conflicts=``) means
# bty-tui starts BEFORE getty has a chance to render its login
# prompt -- otherwise the operator sees the agetty login flash on
# screen for the seconds it takes us to take over, which they read
# as "boot didn't go directly to the TUI". With Before+Conflicts,
# getty never gets a chance to start in the first place.
Conflicts=getty@tty1.service
Before=getty@tty1.service
# Wait for plymouth to fully release tty1 before bty-tui claims
# it. ``plymouth-quit.service`` only SIGNALS plymouth to start
# exiting; ``plymouth-quit-wait.service`` BLOCKS until the
# plymouth process actually terminates. Need both to avoid the
# race that bit second-boot hangs (first boot: grow-images
# takes ~30s so plymouth has time to exit; second boot: grow is
# no-op, plymouth-quit fires immediately, bty-tui starts before
# plymouth releases the tty -> Textual can't initialize ->
# unit hangs forever). The ordering is compatible with
# ``Before=getty@tty1.service`` (plymouth-quit-wait -> bty-tui
# -> getty in the dependency graph).
After=plymouth-quit.service plymouth-quit-wait.service
# ``Wants=network-online.target`` (without a matching ``After=``)
# pulls network up in parallel rather than blocking the TUI on it.
# USB-local mode (``bty-tui-on-tty1`` with no ``bty.server=``) doesn't
# need network at all; PXE-tui mode catalogs over the network but
# the operator typically takes longer to navigate the menu than
# DHCP+linkup take to settle, so racing them is fine. If the
# initial catalog fetch fires before network is up, ``bty-tui``'s
# refresh binding (``r``) retries.
Wants=network-online.target

[Service]
Type=idle
ExecStart=/usr/local/sbin/bty-tui-on-tty1
Restart=on-failure
RestartSec=5
# Plumb tty1 in as the controlling terminal so textual's input loop
# works (cursor keys, function keys, etc.).
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
