# bty-media/live-build/config/includes.chroot/etc/systemd/system/bty-grow-images-partition.service
[Unit]
Description=Grow bty BTY_IMAGES exFAT partition to fill the boot stick
Documentation=https://github.com/safl/bty
# The cooked usb-iso ISO carries a 4 GiB BTY_IMAGES exFAT partition
# pre-allocated at the trailing edge. After dd'ing to a larger stick,
# this oneshot reclaims the unallocated tail (only on first boot,
# only if the partition is empty - see the script for the safety
# gates). No-op on netboot live-x86 (no BTY_IMAGES partition) and
# Ventoy (live medium is a loopback, not a real block device).
#
# Must run before bty-tui-on-tty1.service and bty-flash-on-boot.service
# so the grown partition is in place by the time anything tries to
# scan it for images. local-fs.target is the right ordering anchor;
# we don't actually need any specific filesystem mounted, just
# /run/live/medium populated by live-boot, which is in place by the
# time local-fs runs.

After=local-fs.target
Before=bty-tui-on-tty1.service bty-flash-on-boot.service

# Don't start on a system that doesn't have a BTY_IMAGES partition
# at all. Lazy-evaluated by systemd; covers netboot live-x86 cleanly
# (no /dev/disk/by-label/BTY_IMAGES symlink there).
ConditionPathExists=/dev/disk/by-label/BTY_IMAGES

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/bty-grow-images-partition
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
