# bty-media/live-build/config/includes.chroot/etc/systemd/system/bty-banner-late.service
[Unit]
Description=bty boot banner -- step 3 of 3 (multi-user reached)
Documentation=https://github.com/safl/bty
# Just after multi-user.target. The banner is informational --
# DO NOT add ``Before=bty-tui-on-tty1.service`` or
# ``Before=bty-flash-on-boot.service`` here even though it would
# read nicely: those units are themselves ``WantedBy=multi-user.
# target``, so ordering this unit ``Before=`` them while also
# being ``After=multi-user.target`` creates an ordering cycle.
# systemd silently breaks the cycle by dropping one of the units
# from the boot transaction. Same trap as the server's
# ``bty-banner-late.service`` -- surfaced on v0.22.4 when bty-web
# got dropped on the server appliance.
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/bty-boot-banner late
StandardOutput=journal+console
StandardError=journal+console
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
