# bty-media/live-build/config/includes.chroot/etc/systemd/system/bty-tui-on-tty1.service
[Unit]
Description=bty-tui interactive flash on tty1
Documentation=https://github.com/safl/bty
# Fire whenever ``bty.mode=interactive`` is on the kernel cmdline.
# Two paths set it:
#   - PXE-tui: server's ``ipxe_tui.j2`` (``boot_policy=tui``)
#   - USB-local: the usb-iso live-build bake sets it directly via
#     ``--bootappend-live`` (M19 phase 2).
# ``ipxe_flash.j2`` (``boot_policy=flash``) does NOT set it, so the
# auto-flash boots fall through to ``bty-flash-on-boot.service``.
ConditionKernelCommandLine=bty.mode=interactive
# Take over tty1: the live env's auto-login agetty is the default
# session there; we want to BE the session, not run inside it.
# ``Before=getty@tty1.service`` (combined with ``Conflicts=``) means
# bty-tui starts BEFORE getty has a chance to render its login
# prompt -- otherwise the operator sees the agetty login flash on
# screen for the seconds it takes us to take over, which they read
# as "boot didn't go directly to the TUI". With Before+Conflicts,
# getty never gets a chance to start in the first place.
Conflicts=getty@tty1.service
Before=getty@tty1.service
# ``After=plymouth-quit.service`` waits for plymouth to release the
# tty1 framebuffer before bty-tui claims it. Without this, plymouth
# can still hold the console when systemd starts bty-tui-on-tty1
# and the unit appears to hang at "plymouth-quit.service" because
# Textual can't initialize on a tty plymouth still owns. The
# ordering is compatible with the ``Before=getty@tty1.service``
# directive above (both can be honoured in the same dependency
# graph: plymouth-quit -> bty-tui -> getty).
After=plymouth-quit.service
# ``Wants=network-online.target`` (without a matching ``After=``)
# pulls network up in parallel rather than blocking the TUI on it.
# USB-local mode (``bty-tui-on-tty1`` with no ``bty.server=``) doesn't
# need network at all; PXE-tui mode catalogs over the network but
# the operator typically takes longer to navigate the menu than
# DHCP+linkup take to settle, so racing them is fine. If the
# initial catalog fetch fires before network is up, ``bty-tui``'s
# refresh binding (``r``) retries.
Wants=network-online.target

[Service]
Type=idle
ExecStart=/usr/local/sbin/bty-tui-on-tty1
Restart=on-failure
RestartSec=5
# Plumb tty1 in as the controlling terminal so textual's input loop
# works (cursor keys, function keys, etc.).
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
