# bty-media/live-build/config/includes.chroot/etc/systemd/system/bty-banner-late.service
[Unit]
Description=bty boot banner -- step 3 of 3 (multi-user reached)
Documentation=https://github.com/safl/bty
# Just before bty-tui-on-tty1 / bty-flash-on-boot take tty1.
# If the operator sees step 3 but bty-tui doesn't appear,
# the wedge is in bty-tui itself (slow Textual / Rich init,
# catalog fetch hang, etc.) -- not in the boot chain.
After=multi-user.target
Before=bty-tui-on-tty1.service bty-flash-on-boot.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/bty-boot-banner late
StandardOutput=journal+console
StandardError=journal+console
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
