# bty-media/rootfs/server/etc/systemd/system/bty-banner-late.service
[Unit]
Description=bty boot banner -- step 3 of 3 (multi-user reached)
Documentation=https://github.com/safl/bty
# Just after multi-user.target. The banner is informational --
# DO NOT add ``Before=bty-web.service`` here even though it would
# read nicely: bty-web.service is itself ``WantedBy=multi-user.
# target``, so ordering this unit ``Before=`` it while also being
# ``After=multi-user.target`` creates an ordering cycle. systemd
# silently breaks the cycle by dropping bty-web.service from the
# boot transaction (caught only by ``journalctl | grep 'ordering
# cycle'``); the appliance then comes up with no /healthz and the
# release fails at test-pxe. Surfaced on v0.22.4 -- never reuse
# the Before= relation here.
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/bty-boot-banner late
StandardOutput=journal+console
StandardError=journal+console
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
