# bty-media/rootfs/server/etc/systemd/system/bty-web-init.service
[Unit]
Description=bty-web first-boot initialisation (state dir, /etc/default/bty-web, login banner)
Documentation=https://github.com/safl/bty
# Skip on every boot after the first - the env file has been written
# and the operator may have edited it.
ConditionPathExists=!/etc/default/bty-web
Before=bty-web.service
DefaultDependencies=no
After=local-fs.target
# Also order after the (nofail) /var/lib/bty mount so first-boot dir
# creation + schema init land on the migrated disk if one is present
# (no-op when /var/lib/bty is just the rootfs).
After=var-lib-bty.mount

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/bty-web-init
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
