# bty-media/rootfs/server/etc/systemd/system/bty-ssh-host-keys.service
[Unit]
Description=Regenerate SSH host keys on first boot of the cooked image
Documentation=https://github.com/safl/bty
# Cooked-image security: the bake leaves /etc/ssh/ssh_host_* removed
# (see ``cloudinit-base-server.user`` runcmd / ``rpi_image_customize.py``)
# so the public release artifact does not carry per-image-build host
# keys that every operator's appliance would otherwise share. This
# oneshot regenerates a unique set on first boot of each cooked
# instance.
#
# Idempotent: the condition checks for the ed25519 key (always
# present after ``ssh-keygen -A``); on every boot after the first,
# the unit is a no-op.
ConditionPathExists=!/etc/ssh/ssh_host_ed25519_key
# Must run before ssh.service / ssh.socket so sshd has keys when
# it starts.
Before=ssh.service ssh.socket
After=local-fs.target
DefaultDependencies=no

[Service]
Type=oneshot
# ``-A`` generates one key per missing type (rsa, ecdsa, ed25519).
ExecStart=/usr/bin/ssh-keygen -A
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
