# bty-media/rootfs/server/etc/systemd/system/bty-web.service
[Unit]
Description=bty-web HTTP server
Documentation=https://github.com/safl/bty
Requires=bty-web-init.service
After=network-online.target bty-web-init.service
Wants=network-online.target

[Service]
Type=exec
EnvironmentFile=/etc/default/bty-web
ExecStart=/opt/bty/venv/bin/bty-web
Restart=on-failure
RestartSec=5
User=bty
Group=bty
ProtectSystem=full
ProtectHome=true
# bty-web's only sudo'd helper as of v0.18 is bty-web-tftp, which
# runs ``systemctl <action> dnsmasq.service`` -- a dbus call, no
# filesystem write. No ReadWritePaths punch-throughs needed.
# (Pre-v0.18 had activate/deactivate helpers writing /etc/default/
# and /etc/dnsmasq.d/; both are gone now that bty doesn't manage
# its own DHCP/TFTP daemons.)
# Intentionally NOT setting NoNewPrivileges=true: bty-web's
# /ui/login authenticates via PAM (pamela -> libpam -> pam_unix ->
# setgid unix_chkpwd to read /etc/shadow). NoNewPrivileges blocks
# the setgid privilege gain that pam_unix needs, so every login
# would fail "authentication failure" with the correct password.
# The bty user is already unprivileged; the additional defense-in-
# depth NNP would offer is small, and breaks PAM.

[Install]
WantedBy=multi-user.target
