# bty-media/rootfs/server/etc/systemd/system/bty-web.service
[Unit]
Description=bty-web HTTP server
Documentation=https://github.com/safl/bty
Requires=bty-web-init.service
After=network-online.target bty-web-init.service
Wants=network-online.target

[Service]
Type=exec
EnvironmentFile=/etc/default/bty-web
ExecStart=/opt/bty/venv/bin/bty-web
Restart=on-failure
RestartSec=5
User=bty
Group=bty
ProtectSystem=full
ProtectHome=true
# Intentionally NOT setting NoNewPrivileges=true: bty-web's
# /auth/login authenticates via PAM (pamela -> libpam -> pam_unix
# -> setgid unix_chkpwd to read /etc/shadow). NoNewPrivileges
# blocks the setgid privilege gain that pam_unix needs, so every
# login would fail "authentication failure" with the correct
# password. The bty user is already unprivileged; the additional
# defense-in-depth NNP would offer is small, and breaks PAM.

[Install]
WantedBy=multi-user.target
